# packages/stmpclean/stmpclean.cron
#!/bin/sh
# $Owl: Owl/packages/stmpclean/stmpclean.cron,v 1.4 2010/09/02 21:30:07 solar Exp $

/usr/sbin/stmpclean -t 10d /tmp
if [ "`readlink /var/tmp`" != "/tmp" ]; then
	/usr/sbin/stmpclean -t 10d /var/tmp
fi
if [ -d /var/cache/man ]; then
	/usr/sbin/stmpclean -t 10d /var/cache/man/{,X11R6/,local/}cat[123456789n]
fi
