# packages/stmpclean/stmpclean.cron
#!/bin/sh
# $Id: Owl/packages/stmpclean/stmpclean.cron,v 1.1 2002/03/30 01:13:18 solar Exp $

/usr/sbin/stmpclean -t 10d /tmp /var/tmp

if [ -d /var/catman ]; then
	/usr/sbin/stmpclean -t 10d /var/catman/{,X11R6/,local/}cat[123456789n]
fi
