# packages/stmpclean/stmpclean.cron
#!/bin/sh
# $Id: Owl/packages/stmpclean/stmpclean.cron,v 1.2 2004/11/08 21:12:13 mci Exp $

/usr/sbin/stmpclean -t 10d /tmp /var/tmp

if [ -d /var/cache/man ]; then
	/usr/sbin/stmpclean -t 10d /var/cache/man/{,X11R6/,local/}cat[123456789n]
fi
