# packages/acct/acct.init
#!/bin/sh
# $Id: Owl/packages/acct/acct.init,v 1.4 2002/03/23 21:21:44 solar Exp $
#
# chkconfig: 2345 35 65
# description: \
#	switch process accounting on or off

ACCTON=/sbin/accton

test -x $ACCTON || exit 0

case "$1" in
start)
	echo "Starting process accounting"
	$ACCTON /var/account/pacct
	;;
stop)
	echo "Stopping process accounting"
	$ACCTON
	;;
restart)
	$0 stop
	$0 start
	;;
*)
	echo "Usage: acct {start|stop|restart}"
	exit 1
esac

exit $?
